# FILES/dropbear-user-limit.service
[Unit]
Description=Limit Dropbear user login

[Service]
Type=simple
ExecStart=/usr/bin/dropbear-user-limit

[Install]
WantedBy=multi-user.target